# install/vagrant/resources/supervisor/vagrant/nginx.conf
[program:nginx]
directory = /usr/sbin
command = nginx
autostart = true
autorestart = true
stdout_logfile = /var/log/dontWasteTime/nginx_out.log
stderr_logfile = /var/log/dontWasteTime/nginx_err.log
